# etc/systemd/system/openvpn@.service.d/override.conf
[Service]
LimitNPROC=infinity
Environment="OPENVPN_OPTS=--mute-replay-warnings"